# apps/pip/com.tandem.pip.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for Pip (global hotkey + desktop widget).
  Prefer: tandem pip autostart install
  Manual:
    1. Edit paths if needed, copy to ~/Library/LaunchAgents/com.tandem.pip.plist
    2. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.tandem.pip.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tandem.pip</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/tandem/apps/pip/run.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/tandem-pip.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tandem-pip.err.log</string>
</dict>
</plist>
