# apps/clippy/com.tandem.clippy.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for Tandem Clippy (global hotkey + task widget).
  1. Copy to ~/Library/LaunchAgents/com.tandem.clippy.plist
  2. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.tandem.clippy.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tandem.clippy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/tandem/apps/clippy/run.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/tandem-clippy.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tandem-clippy.err.log</string>
</dict>
</plist>
